# optimizer/docker/crontab
# Crontab for Portfolio Optimizer Daily Analysis
# Runs macro regime and stock signal analysis at midnight UTC

# Set shell
SHELL=/bin/bash

# Set PATH to include Python
PATH=/usr/local/bin:/usr/bin:/bin

# Load environment variables from .env file (if exists)
# Note: Environment variables should be passed via docker-compose or docker run

# Run daily analysis at 22:00 Italy time (Europe/Rome timezone)
# Container timezone is set to Europe/Rome in docker-compose.yml
# This ensures consistent 22:00 local time regardless of DST changes
# Output goes to both cron.log AND Docker logs (via /proc/1/fd/1)
# Format: minute hour day month day_of_week command
0 22 * * * /app/docker/run_daily_analysis.sh 2>&1 | tee /proc/1/fd/1 /var/log/cron.log

# Empty line required at end of crontab file
